# udev.rules
SUBSYSTEM=="hid", ACTION=="add", ENV{HID_ID}=="0003:000028BD:*", \
	RUN+="hid-rebind"
